# serviceTest/ficDownload.service
[Unit]
Description=FicDownloader

[Service]
ExecStart=/home/angel/epub/webTest.py

[Install]
WantedBy=multi-user.target
